# .service
[Unit]
Description=Vivo Webinar Application - v1
After=network.target

[Service]
User={{ application_user }}
Group=webpoint
WorkingDirectory={{ application_path }}/{{ application_user }}/public_html/
ExecStart={{ application_path }}/{{ application_user }}/public_html/env/bin/gunicorn -c gunicorn.conf.py wsgi:app
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target